# deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
autostart=true
autorestart=true
priority=10

[program:nginx]
command=sh -c "nginx -g 'daemon off;'"
autostart=true
autorestart=true
priority=20
environment=PORT="10000"
